# supervisord.conf
[supervisord]
nodaemon=true
# supervisor itself runs as root, because it must be able to open /dev/stdout (the container's
# stdout pipe is root-owned) to relay its children's logs to `docker logs`. It only supervises;
# both programs below drop to the unprivileged app user, so nothing that touches the network
# ever runs as root.
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=dicomium
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true

# Single worker, deliberately: the ingest job registry and its cancel flags live in
# process memory. A second worker would not see them.
[program:backend]
command=/app/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
user=dicomium
directory=/app/backend
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",HOME="/home/dicomium"
